# mpbox/packages/busybox/files/inittab
# Begin /etc/inittab

::sysinit:/etc/sysinit

::respawn:/sbin/udevd --daemon

# End /etc/inittab
